# crontab
0 12 * * * ~/.cronjobs/update.sh >> ~/.cronjobslog.log 2>&1
